# conf/supervise-opencanary.conf
[program:opencanaryd]
command="/opt/opencanary/scripts/startcanary.sh"
redirect_stderr=true
autorestart=true
